# app-misc/ollama-bin/files/ollama.initd
#!/sbin/openrc-run

name=$RC_SVCNAME
description="Start ollama's daemon"
supervisor=supervise-daemon
group=${group:-${RC_SVCNAME}}
user=${user:-${RC_SVCNAME}}
command="/usr/bin/ollama"
command_args="serve"
command_user="ollama"
